# HOME/private_dot_config/systemd/user/redshift.service
[Unit]
Description=RedShift

[Service]
Type=idle
ExecStart=/usr/bin/redshift

[Install]
WantedBy=default.target